# profile.d/__init__.sh
# Ensure that ./node_modules/.bin shows up early in $PATH
export PATH=./node_modules/.bin:$PATH

# Ensure that /usr/local/bin shows up early in $PATH
export PATH=/usr/local/bin:$PATH
