# profile.d/goenv.sh
export GOENV_ROOT=/usr/local/var/goenv
(which -s goenv &> /dev/null) && eval "$(goenv init -)"
