# profile.d/nodenv.sh
export NODENV_ROOT=/usr/local/var/nodenv
(which -s nodenv &> /dev/null) && eval "$(nodenv init -)"
