# profile.d/rbenv.sh
export RBENV_ROOT=/usr/local/var/rbenv
(which -s rbenv &> /dev/null) && eval "$(rbenv init -)"
